# php7/service/PullLocalDaemon/run
#!/usr/bin/env bash
/srv/phabricator/bin/phd launch PhabricatorRepositoryPullLocalDaemon
